# odroid_remote_service.conf
# ODROID Remote Control

description	"ODROID Remote Control"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/bin/odroid_remote /etc/odroid_remote.conf

